# files/pihole/pihole-pull.timer
[Unit]
Description=Run pihole-pull daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
