# files/failure-tg@.service
[Unit]
Description=OnFailure for %i

[Service]
Type=oneshot
ExecStart=/usr/local/bin/onfailure-tg.sh %i
