# files/nginx/certbot.timer
[Unit]
Description=Run certbot weekly and on boot

[Timer]
OnBootSec=70min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
